# systemd/agentd.service
[Unit]
Description=OS-Agent Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/agentd --socket /run/agentd/agentd.sock
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/agentd/env
RuntimeDirectory=agentd
User=root
Group=root

[Install]
WantedBy=multi-user.target
