# systemd/aletheond.service
[Unit]
Description=Aletheon AI Agent Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/aletheond --socket /run/aletheond/aletheond.sock
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/aletheond/env
RuntimeDirectory=aletheond
User=root
Group=root

[Install]
WantedBy=multi-user.target
