# config/aletheon.service
[Unit]
Description=Aletheon AI Agent Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/aletheon daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

RuntimeDirectory=aletheon
RuntimeDirectoryMode=0755
StateDirectory=aletheon
StateDirectoryMode=0755

# Security hardening
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
