# config/aletheon.user.service
[Unit]
Description=Aletheon per-user runtime
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/aletheon daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
RuntimeDirectory=aletheon
RuntimeDirectoryMode=0700
StateDirectory=aletheon
CacheDirectory=aletheon

# The runtime deliberately retains the invoking user's filesystem authority.
# Per-turn WorkspacePolicy and bubblewrap narrow that authority dynamically.
NoNewPrivileges=yes
LimitCORE=0
