# config/aletheon.service
[Unit]
Description=Aletheon AI Agent Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/aletheon daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

RuntimeDirectory=aletheon
RuntimeDirectoryMode=0755
StateDirectory=aletheon
StateDirectoryMode=0755
ReadWritePaths=/etc/aletheon

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ReadWritePaths=%h/.local/share/aletheon

[Install]
WantedBy=multi-user.target
