# config/aletheon-cleanup.service
[Unit]
Description=Policy-expired Aletheon managed-storage cleanup
Documentation=file:/usr/share/doc/aletheon/storage-policy.md
After=aletheon-backup.service

[Service]
Type=oneshot
User=aletheon
Group=aletheon
UMask=0027
ExecStart=/usr/libexec/aletheon/cleanup-aletheon.sh
Environment=ALETHEON_DATA_ROOT=/var/lib/aletheon
Environment=ALETHEON_CACHE_ROOT=/var/cache/aletheon
ReadWritePaths=/var/lib/aletheon /var/cache/aletheon
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
LimitCORE=0
