# config/aletheon.service
[Unit]
Description=Aletheon AI Agent Daemon
After=network.target

[Service]
Type=simple
User=aurobear
ExecStartPre=/usr/bin/mkdir -p /home/aurobear/.local/share/aletheon /home/aurobear/.aletheon /run/aletheon
ExecStart=/usr/bin/aletheon daemon
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/aletheon/.env
Environment=HOME=/home/aurobear
Environment=AGENT_DATA_DIR=/home/aurobear/.local/share/aletheon
Environment=RUST_LOG=info
WorkingDirectory=/home/aurobear/Bear-ws/aletheon

RuntimeDirectory=aletheon
RuntimeDirectoryMode=0755
StateDirectory=aletheon
StateDirectoryMode=0755

# Security hardening
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
